# packaging/systemd/pangolin@.service
# systemd template unit for pangolin.
#
# One instance per saved portal profile. After installing
# (see packaging/systemd/README.md), enable a profile with:
#
#     sudo systemctl enable --now pangolin@work.service
#
# `%i` (instance name) becomes both the pgn instance identifier
# (drives the control socket at /run/pangolin/<instance>.sock) and
# the positional portal argument. Passing it via --instance is what
# lets several `pangolin@*.service` units run in parallel without
# fighting over the same socket.
#
# Instance names are validated by pgn at parse time to
# [A-Za-z0-9_-]{1,32}. Bare URLs are NOT supported here — if you
# want to point a unit at a raw portal URL, save it as a profile
# first with `pgn portal add <name> --url …`.
#
# This is the simplest viable systemd integration. It does NOT
# yet:
#   * Drop privileges below root via AmbientCapabilities (pgn
#     needs CAP_NET_ADMIN for tun + ip route, but the rest of
#     the process can be confined). Phase 2b polish.
#   * Wire ExecStartPre into a portal-availability check.
#   * Use `pgn status --json` for a Type=notify ready signal.
#
# All three are queued as follow-ups.

[Unit]
Description=Pangolin GlobalProtect VPN client (profile %I)
Documentation=https://github.com/kyaky/pangolin
After=network-online.target
Wants=network-online.target

# Auto-restart budget. StartLimitIntervalSec / StartLimitBurst
# belong in [Unit] per systemd.unit(5). Cap how aggressive the
# retries get so a permanently broken portal (e.g. expired
# credentials) gets backed off.
StartLimitIntervalSec=10min
StartLimitBurst=5

[Service]
Type=simple

# Foreground operation — pgn connect blocks until the tunnel
# exits. `--instance %i` pins the control socket to
# /run/pangolin/%i.sock so two parallel units (e.g.
# pangolin@work + pangolin@client-a) never collide.
ExecStart=/usr/local/bin/pgn connect --instance %i %i

# Clean shutdown: systemd sends SIGTERM, pgn's signal handler
# routes that through the libopenconnect cmd pipe and tears the
# tunnel down on the same code path as Ctrl-C. No ExecStop= —
# `pgn disconnect` is asynchronous (it returns as soon as the
# IPC server acknowledges, not when the tunnel is actually down)
# and relying on it here would race the SIGTERM path.
TimeoutStopSec=30s

# Auto-restart on tunnel failure. RestartSec gives the network
# a moment to settle and avoids a tight respawn loop on a
# permanently broken portal.
Restart=on-failure
RestartSec=15s

# Capabilities: pgn manipulates the tun device and the routing
# table. It currently runs as root via the install convention
# (the binary is in /usr/local/bin and the service inherits
# `User=root`). A future hardening pass can switch to
# `User=pangolin` + AmbientCapabilities + NoNewPrivileges, but
# that requires kernel-side audit and is out of scope here.
User=root

# Send stdout/stderr to the journal. pgn already uses tracing
# with structured fields; journald picks them up automatically.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
